# doc/guide/testing.pd
Testing the Installation
========================

TODO

Create the following script for startup, namely in `/etc/init.d/rtai`

    #!/bin/sh

    #------------------------------------------------------------------------------
    #
    # Init script for insmod/rmmoding rtai modules
    #
    #------------------------------------------------------------------------------

    ### BEGIN INIT INFO
    # Provides:          Loading of rtai modules
    # Required-Start:
    # Should-Start:
    # Required-Stop:
    # Should-Stop:
    # Default-Start:     3 5
    # Default-Stop:      0 1 2 6
    # Short-Description: RTAI loader
    # Description:       RTAI loader
    ### END INIT INFO

    COMMAND=start
    if test $# -ge 1
    then
    	COMMAND=$1
    fi

    case $COMMAND in

    start)

        if test \! -c /dev/rtai_shm; then
            mknod -m 666 /dev/rtai_shm c 10 254
        fi
        for n in `seq 0 9`; do
            f=/dev/rtf$n
            if test \! -c $f; then
                mknod -m 666 $f c 150 $n
            fi
        done

        insmod /usr/realtime/modules/rtai_hal.ko
        insmod /usr/realtime/modules/rtai_sched.ko
        insmod /usr/realtime/modules/rtai_sem.ko
        insmod /usr/realtime/modules/rtai_shm.ko
        insmod /usr/realtime/modules/rtai_mbx.ko
        insmod /usr/realtime/modules/rtai_rtdm.ko
        exit 0
        ;;

    stop)

        rmmod rtai_rtdm
        rmmod rtai_mbx
        rmmod rtai_shm
        rmmod rtai_sem
        rmmod rtai_sched
        rmmod rtai_hal
        exit 0
        ;;

    *)
        echo "Usage: $0 {start|stop}"
        ;;

    esac

And make Linux load and unload it automatically:

    $ sudo update-rc.d rtai defaults
